# files/etc/init.d/vpn
#!/bin/sh /etc/rc.common

START=80
STOP=20

restart()
{
    ubus call ipsec reload
}
